# Boxes/SpongeStore/cleanup.timer
[Unit]
Description=Cleanup
Requires=cleanup.service

[Timer]
Unit=cleanup.service
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target